# scripts/com.codexmonitor.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.codexmonitor.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABS/PATH/TO/codex_monitor_daemon</string>
    <string>--listen</string><string>127.0.0.1:4732</string>
    <string>--data-dir</string><string>/Users/YOU/Library/Application Support/codex-monitor-daemon</string>
    <string>--token</string><string>REPLACE_WITH_STRONG_TOKEN</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/codex-monitor-daemon.out.log</string>
  <key>StandardErrorPath</key><string>/tmp/codex-monitor-daemon.err.log</string>
</dict>
</plist>
